# certbot.cron
17 22 * * * root certbot renew --renew-hook "/etc/init.d/nginx reload" > /dev/stdout
