# certbot.cron
17 05,35 * * * root certbot renew --renew-hook "/etc/init.d/nginx reload" > /dev/stdout
